# dot_config/systemd/user/borg-backup.timer
[Unit]
Description=Run Backup and rclone Daily

[Timer]
OnCalendar= *-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
