# dot_config/systemd/user/emacs-agenda.timer
[Unit]
Description=Run the agenda export service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
